# com.pickwinningnumbers.results-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.pickwinningnumbers.results-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>'/Users/ianstinson/Library/CloudStorage/GoogleDrive-signup06@gmail.com/My Drive/aWebSites/PickWinningNumbers/scripts/update_results.sh'</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ianstinson/Library/CloudStorage/GoogleDrive-signup06@gmail.com/My Drive/aWebSites/PickWinningNumbers</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ianstinson/Library/CloudStorage/GoogleDrive-signup06@gmail.com/My Drive/aWebSites/PickWinningNumbers/results-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ianstinson/Library/CloudStorage/GoogleDrive-signup06@gmail.com/My Drive/aWebSites/PickWinningNumbers/results-refresh-error.log</string>
  </dict>
</plist>
